# templates/arangodbcluster.service
[Unit]
Description=Arangodb Cluster

[Service]
ExecStart=/usr/bin/arangodb {{option}}
